# dnsmasq.service.d/all-servers.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dnsmasq -k --enable-dbus --user=dnsmasq --pid-file --all-servers
